# .etc/crontab/ipmonitor
* * * * * cd /opt/ipmonitor && sudo -E docker compose exec workspace php artisan schedule:run >> /dev/null 2>&1
